# systemd/user/pacman-updates.service
[Unit]
Description=Check for pacman and AUR updates

[Service]
Type=oneshot
ExecStart=%h/.local/bin/statusbar/slstatus-upd

